# docker_app/sv_uwsgi.conf
[program:uwsgi]
command=/var/env/bin/uwsgi --ini /etc/uwsgi.ini
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=false
