# docker_app/sv/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/www/app/log/%(program_name)s.log
stderr_logfile=/var/www/app/log/%(program_name)s.log
autorestart=true