# docker_app/sv/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
stdout_logfile=/var/www/app/log/%(program_name)s.log
stderr_logfile=/var/www/app/log/%(program_name)s.log
autorestart=true
